# crontab
00 19 * * * root python main.py >> sts.log
