# crontab
0 18 * * * python bot.py >> taleb-bot.log
